# conf.d/ss-supervisord.conf
[program:ss-redir]
command=/usr/local/bin/ss-redir -v -u -s {ip} -p {port} -k {pass} -m {method} -l 1080 --fast-open -b 0.0.0.0
autostart=true
autorestart=unexpected
stdout_logfile=/tmp/ss-server.log

[program:ss-tunnel]
command=/usr/local/bin/ss-tunnel -v -s {ip} -p {port} -k {pass} -m {method} -L 8.8.8.8:53 -l 5353 -u -b 0.0.0.0
autostart=true
autorestart=unexpected
stdout_logfile=/tmp/ss-tunnel.log
